# files/core/etc/systemd/system/nas-media.automount
[Automount]
Where=/mnt/nas/Medien
TimeoutIdleSec=1min

[Install]
WantedBy=multi-user.target
